# docker/ros_ws_src/ainex_bringup/service/start_app_node.service
[Unit]
Description=start node
After=NetworkManager.service time-sync.target

[Service]
Type=simple
Restart=always
RestartSec=30
KillMode=mixed
ExecStart=/home/ubuntu/ros_ws/src/ainex_bringup/scripts/source_env.bash roslaunch ainex_bringup bringup.launch
StandardOutput=null
StandardError=null

[Install]
WantedBy=multi-user.target
